# rootfs/etc/systemd/system/quickfw-api.service
[Unit]
Description=QuickFW API Server & Web Dashboard
After=network-online.target nftables.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
Environment=STATIC_FILES_PATH=/opt/quickfw/front
ExecStart=/usr/local/bin/quickfw-api --log-level info
LimitNOFILE=65536
Restart=on-failure
RestartSec=3s
StartLimitIntervalSec=300
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quickfw-api

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_KILL CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
